# shell/backup.service
[Unit]
Description= Backup btrfs subvolumes to external disk

[Service]
ExecStart=/home/flo/.local/bin/backup
